# src/crontab
10 2 * * * . /env.sh && /backup.sh > /proc/1/fd/1 2>&1
