# src/crontab
10 2-4 */4 1,2,*/3 * source /env.sh && rcb backup > /proc/1/fd/1

