# src/crontab
10 2 * * * source /env.sh && /backup.sh && rcb cleanup > /proc/1/fd/1

