# src/crontab
10 2 * * * source /env.sh && /backup.sh > /proc/1/fd/1

